# server/systemd/fleet-derp-redirect.service
[Unit]
# 仅 NAT 场景需要（对外端口 ≠ 本机监听端口）：网关在 NAT 后、做了 公网__HS_PUBLIC_PORT__→本机__HS_LISTEN_PORT__
# 端口转换且 hairpin 不通时，让网关自身访问 127.0.0.1:__HS_PUBLIC_PORT__（embedded DERP 对外端口）
# 重定向到本地 __HS_LISTEN_PORT__，这样网关 tailscaled 也能连上自己的 DERP，经中继到达各 Mac。
# 占位 __HS_PUBLIC_PORT__/__HS_LISTEN_PORT__ 由 setup-server.sh 按 .env 渲染，且仅在两端口不同时安装启用。
Description=Loopback redirect for gateway embedded DERP (NAT hairpin workaround)
After=network.target
Before=tailscaled.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/iptables -t nat -A OUTPUT -p tcp -d 127.0.0.1 --dport __HS_PUBLIC_PORT__ -j REDIRECT --to-ports __HS_LISTEN_PORT__
ExecStop=/usr/sbin/iptables -t nat -D OUTPUT -p tcp -d 127.0.0.1 --dport __HS_PUBLIC_PORT__ -j REDIRECT --to-ports __HS_LISTEN_PORT__

[Install]
WantedBy=multi-user.target
