# server/systemd/headscale.service
[Unit]
Description=headscale self-hosted Tailscale control server (mac-fleet-hub)
After=network.target
Documentation=https://headscale.net

[Service]
Type=simple
# 以 root 运行以便读取 TLS 私钥（acme.sh 签发，默认 /root/.acme.sh/*，root-only，且会随续期轮换）。
User=root
Group=root
ExecStart=/usr/bin/headscale serve
Restart=always
RestartSec=5
ReadWritePaths=/var/lib/headscale /var/run/headscale
RuntimeDirectory=headscale
StateDirectory=headscale
NoNewPrivileges=true
ProtectSystem=strict
# 证书在 /root/.acme.sh（home 目录）→ 用 read-only 而非 true，否则 ProtectHome 屏蔽 /root 致读不到证书。
# 若改用 acme.sh --install-cert 部署到 /root 之外的固定路径，可改回 ProtectHome=true。
ProtectHome=read-only
PrivateTmp=true
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
