# mac/com.macfleet.ttyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd 服务：ttyd → tmux 持久会话(跑 claude)。
     由 setup-mac.sh 把 __BREW_PREFIX__ / __TS_IP__ / __PORT__ 替换为实际值后，
     安装到 ~/Library/LaunchAgents/ 并 launchctl load。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macfleet.ttyd</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BREW_PREFIX__/bin/ttyd</string>
        <string>-i</string>
        <string>__TS_IP__</string>
        <string>-p</string>
        <string>__PORT__</string>
        <string>-W</string>
        <string>__BREW_PREFIX__/bin/tmux</string>
        <string>new</string>
        <string>-A</string>
        <string>-s</string>
        <string>claude</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/macfleet-ttyd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/macfleet-ttyd.err</string>
    <!-- 让 claude 等命令可被找到 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__BREW_PREFIX__/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
