# mac/com.macfleet.ttyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd 服务：ttyd 网页终端。
     ttyd 接到 tmux（经 fleet-attach.sh，按 URL ?arg=<会话名> 选窗）；会话由 fleet-agent 创建。
     由 setup-mac.sh 替换 __BREW_PREFIX__ / __TS_IP__ / __PORT__ / __TTYD_BASE__ / __FLEET_ATTACH__。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macfleet.ttyd</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BREW_PREFIX__/bin/ttyd</string>
        <string>-i</string>
        <string>__TS_IP__</string>
        <string>-p</string>
        <string>__PORT__</string>
        <string>-b</string>
        <string>__TTYD_BASE__</string>
        <string>-W</string>
        <string>-t</string>
        <string>disableLeaveAlert=true</string>
        <string>--url-arg</string>
        <string>__FLEET_ATTACH__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/macfleet-ttyd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/macfleet-ttyd.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- UTF-8 locale：launchd 不继承用户 locale，缺省下 tmux attach 客户端按非 UTF-8
             处理，会吃掉中文等多字节输入。显式设 UTF-8 才能在网页终端输入中文。 -->
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
        <key>PATH</key>
        <string>__BREW_PREFIX__/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
